# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=5MB
logfile_backups=3
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.err
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.err
priority=20

[program:bore]
command=/usr/local/bin/bore-setup.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/bore.log
stderr_logfile=/var/log/bore.err
priority=40

[program:watchdog]
command=/usr/local/bin/watchdog.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/watchdog.log
stderr_logfile=/var/log/watchdog.err
priority=99
